# {{cookiecutter.out_dir}}/sys/crontab
MAILTO=""
# Launching a cron scheduler
* * * * *  if [ ! -e {{ cookiecutter.private_path }}/SUSPEND_CRONS ];then /code/app/bin/console cron:run;fi
